# src/welle-gui/raspi/welle-io.service
[Unit]
Description=welle-io DAB+ receiver
After=network.target

[Service]
ExecStart=/opt/welle-io/welle-io


[Install]
WantedBy=multi-user.target
